# stages/ro/varlib.service
[Unit]
Description=Creates /.varlib/{upper,work}

[Service]
Type=oneshot
ExecStart=/bin/bash -c "mkdir /.varlib/{upper,work}"
RemainAfterExit=true
